# deploy/speech-to-text.service
[Unit]
Description=Speech to Text Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/jirka/speech-to-text
Environment="DOTNET_ROOT=/home/jirka/.dotnet"
ExecStart=/home/jirka/.dotnet/dotnet /home/jirka/speech-to-text/Olbrasoft.SpeechToText.Service.dll
# No Restart directive - service will not restart automatically when quit

[Install]
WantedBy=default.target
